# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:proxy]
command=/opt/app/bin/proxyd
